# deploy/systemd/ai-web-fetcher.service
[Unit]
Description=AI Web Fetcher — web.{search,fetch}.v1 Unix-socket RPC server
Documentation=https://github.com/hamlet-archer/ai-web-fetcher
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
# Type=simple rather than Type=notify: Node 20 LTS does not expose
# AF_UNIX dgram sockets, so a stock-Node sd_notify(READY=1) is not
# implementable without a native addon. Pragmatic mirror of
# ai-calendar-adviser.service + ai-comms-adviser-query.service.
# Readiness is confirmed by the ExecStartPost probe below.
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/ai-web-fetcher
EnvironmentFile=/opt/ai-web-fetcher/.env
ExecStart=/usr/bin/node dist/main.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
TimeoutStartSec=30
TimeoutStopSec=10

# Block until the RPC socket is actually bound — the unit isn't "active"
# until callers can connect. If main.ts fatals during runBootCheck
# (allowlist DNS drift, missing BRAVE_SEARCH_API_KEY), this fires
# before Restart kicks in and gives systemd a clean failed-state.
ExecStartPost=/bin/bash -c 'for i in $(seq 1 30); do [ -S /run/ai-web-fetcher/query.sock ] && exit 0; sleep 1; done; exit 1'

# RuntimeDirectory creates /run/ai-web-fetcher/ owned by ubuntu:ubuntu
# at service start; systemd cleans it up on stop. main.ts binds its
# socket inside this dir at /run/ai-web-fetcher/query.sock.
RuntimeDirectory=ai-web-fetcher
RuntimeDirectoryMode=0755

# Brave Search API key handed in via LoadCredential. main.ts resolves
# it via $CREDENTIALS_DIRECTORY/brave-api-key (see loadBraveApiKey
# in src/main.ts). The .env may also set BRAVE_SEARCH_API_KEY
# directly for local-dev parity, but in production LoadCredential is
# the chain of custody.
LoadCredential=brave-api-key:/etc/ai-web-fetcher/brave-api-key

# Hardening (mirrors ai-calendar-adviser).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/ai-web-fetcher
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# UDS-only for callers; outbound HTTPS to Brave Search + allowlisted
# domains needs AF_INET + AF_INET6.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
